# etc/s6-overlay/s6-rc.d/tail-log/run
#!/command/with-contenv sh
tail -F /var/www/html/log/errors.log /var/www/html/log/debug.log /var/www/html/log/info.log
