# etc/s6-overlay/s6-rc.d/movim-migrations/up
#!/command/with-contenv sh
composer -d /var/www/html/ movim:migrate