# etc/s6-overlay/s6-rc.d/movim-daemon/run
#!/command/with-contenv sh
php /var/www/html/daemon.php start